# build/deb/stable/DEBIAN/postinst
#!/bin/sh

# start on boot
update-rc.d opensemanticsearch-filemonitoring defaults
